# toonabi.service
[Unit]
Description=Toonabi Rails API
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/var/www/toonabi/backend
Environment=RAILS_ENV=production
Environment=PORT=3001
ExecStart=/home/deploy/.rbenv/shims/bundle exec puma -C config/puma.rb
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
